# roles/munerator_new/templates/munerator.conf
[group:munerator]
programs: {{modules|join(',')}}

{% for module in modules %}
[program:{{module}}]
user=openarena
command=nice -n5 munerator {{module}}
{% endfor %}